# services/etc/systemd/system/ota_update.service
[Unit]
Description=OTA Update Service for Jetson Nano
After=network.target

[Service]
ExecStart=/usr/bin/python3 /opt/updater.py
Restart=on-failure

[Install]
WantedBy=multi-user.target
